# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx-agent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

if bashio::config.has_value 'data_plane_key'; then
    TOKEN=$(bashio::config 'data_plane_key')
    bashio::log.info "Found data_plane_key, applying safe configuration updates..."
    
    # 1. 
    # 
    if grep -q "dummy_key_for_build" /etc/nginx-agent/nginx-agent.conf; then
        sed -i "s#dummy_key_for_build#${TOKEN}#g" /etc/nginx-agent/nginx-agent.conf
        bashio::log.info "Token successfully injected into original config."
    fi

    # 2. 
    # 
    if ! grep -q "allowed_directories:" /etc/nginx-agent/nginx-agent.conf; then
        bashio::log.info "Appending exhaustive directory whitelist..."
        cat <<EOF >> /etc/nginx-agent/nginx-agent.conf

allowed_directories:
  - "/etc"
  - "/usr"
  - "/var"
  - "/data"
  - "/config"
  - "/ssl"
  - "/proc"
EOF
    fi

    bashio::log.info "Starting NGINX Agent process..."
    # 
    exec /usr/bin/nginx-agent
else
    bashio::log.warning "No data_plane_key provided."
    exec sleep infinity
fi
