# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx-agent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Check if the user has provided the Token in the HA UI
if bashio::config.has_value 'data_plane_key'; then
    TOKEN=$(bashio::config 'data_plane_key')
    bashio::log.info "Found data_plane_key, preparing F5 NGINX One Console connection..."
    
    # Create the directory for the Agent configuration file
    mkdir -p /etc/nginx-agent
    
    # Dynamically generate the NGINX Agent configuration file (including log path override)
    cat <<EOF > /etc/nginx-agent/nginx-agent.conf
server:
  host: "agent.connect.nginx.com"
  port: 443
auth:
  token: "${TOKEN}"
tls:
  skip_verify: false

# Enable metrics sender feature
features:
  - metrics-sender

metrics:
  bulk_size: 20
  report_interval: 1m
  collection_interval: 15s
  metrics_poll_interval: 1000ms

# Force the Agent to read the physical Access Log we created (bypassing NPM's redirection of logs to stdout)
nginx:
  access_logs:
    - "/tmp/nginx-agent-access.log"
EOF

    bashio::log.info "NGINX Agent configuration generated successfully."
    bashio::log.info "Starting NGINX Agent process..."
    
    # Use exec to start the main process so s6-overlay can properly manage the PID
    exec /usr/bin/nginx-agent
else
    bashio::log.warning "No data_plane_key provided in the add-on configuration."
    bashio::log.warning "NGINX Agent will enter sleep mode and will not affect NGINX Proxy Manager operations."
    
    # Sleep infinitely to keep the s6-overlay service alive without consuming resources
    exec sleep infinity
fi
