# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx-agent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Check if the user has provided the Token in the HA UI
if bashio::config.has_value 'data_plane_key'; then
    TOKEN=$(bashio::config 'data_plane_key')
    bashio::log.info "Found data_plane_key, preparing F5 NGINX One Console connection..."
    
    # Strategy change: Instead of manually generating YAML, directly replace 
    # the dummy key in the official config generated during the Docker build!
    if grep -q "dummy_key_for_build" /etc/nginx-agent/nginx-agent.conf 2>/dev/null; then
        sed -i "s/dummy_key_for_build/${TOKEN}/g" /etc/nginx-agent/nginx-agent.conf
        bashio::log.info "NGINX Agent token injected into static configuration."
    fi
    
    # Some Agent versions write the Token in the dynamic config, check and replace it there as well.
    if [ -f /var/lib/nginx-agent/agent-dynamic.conf ] && grep -q "dummy_key_for_build" /var/lib/nginx-agent/agent-dynamic.conf 2>/dev/null; then
        sed -i "s/dummy_key_for_build/${TOKEN}/g" /var/lib/nginx-agent/agent-dynamic.conf
        bashio::log.info "NGINX Agent token injected into dynamic configuration."
    fi
    
    # Ultimate fallback: Use officially supported environment variables to override the config.
    export NGINX_AGENT_SERVER_TOKEN="${TOKEN}"

    #
    # Expand the Agent's trusted directories to include NPM's dynamic configuration and SSL folders
    export NGINX_AGENT_CONFIG_DIRS="/etc/nginx:/data/nginx:/etc/letsencrypt:/data/custom_ssl"
    bashio::log.info "NGINX Agent trusted directories expanded for NPM compatibility."

    bashio::log.info "Starting NGINX Agent process..."
    
    # Use exec to start the main process so s6-overlay can properly manage the PID
    exec /usr/bin/nginx-agent
else
    bashio::log.warning "No data_plane_key provided in the add-on configuration."
    bashio::log.warning "NGINX Agent will enter sleep mode and will not affect NGINX Proxy Manager operations."
    
    # Sleep infinitely to keep the s6-overlay service alive without consuming resources
    exec sleep infinity
fi
