# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx-agent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Check if the user has provided the Token in the HA UI
if bashio::config.has_value 'data_plane_key'; then
    TOKEN=$(bashio::config 'data_plane_key')
    bashio::log.info "Found data_plane_key, preparing F5 NGINX One Console connection..."
    
    # 1. Inject the Token into the official config
    # Use '#' as sed delimiter instead of '/' to handle slashes in the Token
    if grep -q "dummy_key_for_build" /etc/nginx-agent/nginx-agent.conf 2>/dev/null; then
        sed -i "s#dummy_key_for_build#${TOKEN}#g" /etc/nginx-agent/nginx-agent.conf
        bashio::log.info "NGINX Agent token injected into static configuration."
    fi
    
    # Check dynamic config as well
    if [ -f /var/lib/nginx-agent/agent-dynamic.conf ] && grep -q "dummy_key_for_build" /var/lib/nginx-agent/agent-dynamic.conf 2>/dev/null; then
        sed -i "s#dummy_key_for_build#${TOKEN}#g" /var/lib/nginx-agent/agent-dynamic.conf
        bashio::log.info "NGINX Agent token injected into dynamic configuration."
    fi

    # 2. Cleanup previous attempts (in case of a simple restart without rebuild)
    # Using '#' here as well for consistency
    sed -i '\#nginx:#,$d' /etc/nginx-agent/nginx-agent.conf
    sed -i '\#allowed_directories:#,$d' /etc/nginx-agent/nginx-agent.conf

    # 3. The Ultimate Fix: Exhaustive exact-match directories!
    cat <<EOF >> /etc/nginx-agent/nginx-agent.conf

allowed_directories:
  - "/config"
  - "/etc/nginx-agent"
  - "/etc/nginx"
  - "/etc/app_protect"
  - "/usr/local/etc/nginx"
  - "/usr/share/nginx/modules"
  - "/var/run/nginx"
  - "/var/log/nginx"
  - "/data/nginx"
  - "/data/custom_ssl"
  - "/etc/letsencrypt"
  - "/config/nginx"
  - "/config/nginx/proxy_host"
  - "/config/nginx/custom"
  - "/config/nginx/redirection_host"
  - "/config/nginx/dead_host"
  - "/config/nginx/stream"
  - "/config/nginx/default_host"
  - "/config/nginx/temp"
  - "/ssl"
  - "/proc"
EOF
    bashio::log.info "NGINX Agent exhaustive trusted directories appended to YAML."

    bashio::log.info "Starting NGINX Agent process..."
    
    # Use exec to start the main process so s6-overlay can properly manage the PID
    exec /usr/bin/nginx-agent
else
    bashio::log.warning "No data_plane_key provided in the add-on configuration."
    bashio::log.warning "NGINX Agent will enter sleep mode and will not affect NGINX Proxy Manager operations."
    
    # Sleep infinitely to keep the s6-overlay service alive without consuming resources
    exec sleep infinity
fi
