# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx-agent/run
#!/command/with-contenv bashio
# shellcheck shell=bash

if bashio::config.has_value 'data_plane_key'; then
    TOKEN=$(bashio::config 'data_plane_key')
    bashio::log.info "Reconstructing NGINX Agent configuration from scratch..."
    
    # 
    mkdir -p /etc/nginx-agent

    #
    cat <<EOF > /etc/nginx-agent/nginx-agent.conf
# Fixed Configuration for NGINX Agent in Home Assistant
log:
  level: info
  path: /var/log/nginx-agent/

allowed_directories:
  - "/etc"
  - "/etc/nginx"
  - "/etc/nginx-agent"
  - "/var/run/nginx"
  - "/var/log/nginx"
  - "/data"
  - "/data/nginx"
  - "/config"
  - "/config/nginx"
  - "/config/nginx/proxy_host"
  - "/config/nginx/custom"
  - "/config/nginx/default_host"
  - "/config/nginx/dead_host"
  - "/config/nginx/redirection_host"
  - "/config/nginx/temp"
  - "/config/nginx/stream"
  - "/ssl"
  - "/proc/1/fd"

command:
  server:
    host: agent.connect.nginx.com
    port: 443
  auth:
    token: "${TOKEN}"
  tls:
    enable: true
    skip_verify: false

features:
  - configuration
  - metrics
  - file-watcher
  - certificates
EOF

    # 
    chmod 644 /etc/nginx-agent/nginx-agent.conf
    
    bashio::log.info "NGINX Agent configuration is now perfectly aligned."
    bashio::log.info "Starting NGINX Agent process..."
    
    exec /usr/bin/nginx-agent
else
    bashio::log.error "No data_plane_key provided!"
    exec sleep infinity
fi
